# monitoring/supervisord/monitoring_uwsgi.ini
[program:uwsgi]
command=/usr/local/webserver/monitor-tel/bin/uwsgi --ini /usr/local/webserver/monitor-tel/mysite/monitoring/conf/monitor-tel_uwsgi.ini
directory=/usr/local/webserver/monitor-tel
autostart=true
autorestart=true
stdout_logfile=/usr/local/webserver/monitor-tel/logs/uwsgi-emperor.log
stderr_logfile=/usr/local/webserver/monitor-tel/logs/uwsgi-emperor.log
stdout_logfile_maxbytes = 20MB
